# ansible/roles/cert-renewal/templates/k8s-cert-renew.service.j2
[Unit]
Description=Kubernetes certificate renewal
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/k8s-cert-renew.sh
StandardOutput=journal
StandardError=journal
